# z_celery/conf/supervisor_celery_beat.ini
[program:imooc-celery-beat]

command=python3 manage.py celery beat -l INFO

directory=/home/tarena/桌面/test/my_code/month04/celery/taskproj

environment=PATH="/home/tarena/桌面/test/my_code/month04/celery/taskproj/venv/bin"

stdout_logfile=/home/tarena/桌面/test/项目实战/z_celery/logs/celery.beat.log

stderr_logfile=/home/tarena/桌面/test/项目实战/z_celery/logs/celery.beat.log

autostart=true

autorestart=true

startsecs=10

stopwatisecs=60

priority=988